# package/base-files/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
#chmod +x /etc/switch.sh
#.//etc/switch.sh &

# We need to delete shadow.fifo pipe in every reboot
# Becasue after reboot this file will be normal file insted of fifo
rm /var/run/shadow.fifo

exit 0
